# cron/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/5 * * * * /usr/local/bin/wp-cron.sh
0 0 * * 0 /optimize-db.sh
